# scripts/clawnsole-updater.launchagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.clawnsole-updater</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>bash</string>
    <string>-lc</string>
    <string>bash "$HOME/.openclaw/apps/clawnsole/scripts/update.sh"</string>
  </array>
  <key>StartInterval</key>
  <integer>21600</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>$HOME/.openclaw/logs/clawnsole.update.out.log</string>
  <key>StandardErrorPath</key>
  <string>$HOME/.openclaw/logs/clawnsole.update.err.log</string>
</dict>
</plist>
